# data/udev/90-HOTPLUG_display-visor.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/darknmt/.Xauthority", RUN+="/usr/bin/pkill -x -RTMIN+5 display-visor"
